# com.ig.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ig.monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/huangke/Developer/ig_scraper/monitor</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/huangke/Developer/ig_scraper</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/huangke/Developer/ig_scraper/monitor.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/huangke/Developer/ig_scraper/monitor.error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
